# jsk_tools/scripts/upstart/robots/staro/roscore.conf
description "roscore daemon"

start on runlevel [2345]
stop on runlevel [!2345]
umask 022

script
     . /opt/ros/hydro/setup.sh
     roscore
end script
